# docker/php/prod/crontab
* * * * * php /var/www/bin/console app:test --no-debug --no-interaction
